# http-master-upstart.conf
author "Rob Calcroft"
description "Start, stop and restart http-master from a service"

start on runlevel [2345]
stop on runlevel [!2345]

# Allow restart on crash
respawn

# Exec
exec http-master --config=/etc/http-master/http-master.conf
